# icecast/play-silence-64-mono-ogg.conf
# Broadcasts silence on icecast

start on stopped rc RUNLEVEL=4

respawn
respawn limit 10 120

exec /usr/local/bin/play-silence-64-mono-ogg
